# tools/configs/brain-overnight.service
# Systemd service file for brain overnight runner
# Install: sudo cp brain-overnight.service /etc/systemd/system/
# Enable: sudo systemctl enable brain-overnight.timer
# Start: sudo systemctl start brain-overnight.timer

[Unit]
Description=Brain System Overnight Runner
After=network.target

[Service]
Type=oneshot
User=%i
Environment="BRAIN_PATH=/home/%i/brain"
ExecStart=/home/%i/brain/tools/configs/overnight-brain.sh
StandardOutput=append:/home/%i/brain/logs/systemd-overnight.log
StandardError=append:/home/%i/brain/logs/systemd-overnight.log
TimeoutStartSec=14400

[Install]
WantedBy=multi-user.target
